# systemd/reader-update-daemon.service
[Unit]
Description=Gheop Reader Feed Update Daemon
After=network.target mysql.service

[Service]
Type=simple
User=sib
Group=sib
WorkingDirectory=/www/reader
ExecStart=/usr/bin/php /www/reader/bin/feeds/update_daemon.php
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=reader-update-daemon

# Security
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
CPUQuota=50%
MemoryLimit=512M

[Install]
WantedBy=multi-user.target
